# apple.conf
[program:apple]
# the path of apple project
directory = '/home/kevin/apple/'

# command must has virtualenv path
command = /home/kevin/py2env/bin/python /home/kevin/apple/apple/web/__init__.py
